# supervisor/conf.d/minecraft.conf
[program:minecraft]
directory=/data/
command=java -Xmx2048M -Xms2048M -jar minecraft_server.jar

